# rootfs/rootfs/etc/services.d/zoho-calendar/run
#!/command/with-contenv bashio
# shellcheck shell=bash

exec /opt/zoho-calendar/run.sh
